# src/upstart/ceph-osd-all.conf
description "Ceph OSD (all instances)"

start on (filesystem and net-device-up)
stop on runlevel [!2345]
