# src/upstart/ceph-all.conf
description "Ceph"

start on runlevel [2345] and static-network-up
stop on runlevel [!2345]
