# src/upstart/ceph-mgr-all.conf
description "Ceph Manager (all instances)"

start on starting ceph-all
stop on stopping ceph-all
