# src/upstart/ceph-mds-all.conf
description "Ceph MDS (all instances)"

start on starting ceph-all
stop on stopping ceph-all
