# src/upstart/ceph-disk.conf
description "ceph-disk async worker"

start on ceph-disk

instance $dev/$pid
export dev
export pid

exec flock /var/lock/ceph-disk -c 'ceph-disk --verbose --log-stdout trigger --sync $dev'

