# src/upstart/radosgw-all.conf
description "Ceph radosgw (all instances)"

start on starting ceph-all
stop on runlevel [!2345] or stopping ceph-all
