# systemd/ceph-fuse@.service.in
[Unit]
Description=Ceph FUSE client
After=network-online.target local-fs.target time-sync.target
Wants=network-online.target local-fs.target time-sync.target
Conflicts=umount.target
PartOf=ceph-fuse.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
Environment=CLUSTER=ceph
ExecStart=/usr/bin/ceph-fuse -f --cluster ${CLUSTER} %I
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
# ceph-fuse requires access to /dev fuse device
PrivateDevices=no
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
TasksMax=infinity
Restart=on-failure
StartLimitInterval=30min
StartLimitBurst=3

[Install]
WantedBy=ceph-fuse.target
