# systemd/ceph-radosgw@.service.in
[Unit]
Description=Ceph rados gateway
PartOf=ceph-radosgw.target
After=network-online.target local-fs.target time-sync.target
Before=remote-fs-pre.target ceph-radosgw.target
Wants=network-online.target local-fs.target time-sync.target remote-fs-pre.target ceph-radosgw.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
ExecStart=/usr/bin/radosgw -f --name client.%i --setuser ceph --setgroup ceph
LimitNOFILE=1048576
LimitNPROC=1048576
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=yes
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full
Restart=on-failure
RestrictSUIDSGID=true
StartLimitBurst=5
StartLimitInterval=30s
TasksMax=infinity

[Install]
WantedBy=ceph-radosgw.target
