# systemd/ceph-crash.service.in
[Unit]
Description=Ceph crash dump collector

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/ceph-crash
Restart=always
RestartSec=10
StartLimitInterval=10min
StartLimitBurst=10

[Install]
WantedBy=ceph.target
