# systemd/ceph-mds@.service.in
[Unit]
Description=Ceph metadata server daemon
PartOf=ceph-mds.target
After=network-online.target local-fs.target time-sync.target
Before=remote-fs-pre.target ceph-mds.target
Wants=network-online.target local-fs.target time-sync.target remote-fs-pre.target ceph-mds.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/ceph-mds -f --id %i --setuser ceph --setgroup ceph
LimitNOFILE=1048576
LimitNPROC=1048576
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=yes
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full
Restart=on-failure
RestrictSUIDSGID=true
StartLimitBurst=3
StartLimitInterval=30min
TasksMax=infinity

[Install]
WantedBy=ceph-mds.target
