# systemd/ceph-fuse@.service.in
[Unit]
Description=Ceph FUSE client
After=network-online.target local-fs.target time-sync.target
Wants=network-online.target local-fs.target time-sync.target
Conflicts=umount.target
PartOf=ceph-fuse.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
ExecStart=/usr/bin/ceph-fuse -f %I
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
# ceph-fuse requires access to /dev fuse device
PrivateDevices=no
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
Restart=on-failure
RestrictSUIDSGID=true
StartLimitBurst=3
StartLimitInterval=30min
TasksMax=infinity

[Install]
WantedBy=ceph-fuse.target
