# systemd/ceph-immutable-object-cache@.service.in
[Unit]
Description=Ceph immutable object cache daemon
After=network-online.target local-fs.target
Wants=network-online.target local-fs.target
PartOf=ceph-immutable-object-cache.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/ceph-immutable-object-cache -f --id %i --setuser ceph --setgroup ceph
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=yes
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full
Restart=on-failure
RestrictSUIDSGID=true
StartLimitBurst=3

[Install]
WantedBy=ceph-immutable-object-cache.target
