# systemd/rbdmap.service.in
[Unit]
Description=Map RBD devices
After=network-online.target ceph.target
Before=remote-fs-pre.target
Wants=network-online.target remote-fs-pre.target ceph.target

[Service]
EnvironmentFile=-@SYSTEMD_ENV_FILE@
Environment=RBDMAPFILE=/etc/ceph/rbdmap
Type=oneshot
RemainAfterExit=yes
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/rbdmap map
ExecReload=@CMAKE_INSTALL_PREFIX@/bin/rbdmap map
ExecStop=@CMAKE_INSTALL_PREFIX@/bin/rbdmap unmap-all

[Install]
WantedBy=multi-user.target
