# systemd/rbdmap-pre.target
[Unit]
Description=Map RBD devices (Pre)
RefuseManualStart=yes
Before=rbdmap.target
