# debian/ceph-mds.prerm
#!/bin/sh
# vim: set noet ts=8:

set -e

invoke-rc.d ceph-mds-all stop || {
	RESULT=$?
	# Ignore if ceph-all upstart config does not
	# exist or upstart is not in use
	if [ $RESULT != 100 ]; then
		exit $RESULT
	fi
}

#DEBHELPER#

exit 0
