# debian/radosgw.prerm
#!/bin/sh
# vim: set noet ts=8:

set -e

case "$1" in
    remove)
	invoke-rc.d radosgw stop || {
	    RESULT=$?
	    if [ $RESULT != 100 ]; then
		exit $RESULT
	    fi
	}
	;;

    *)
	;;
esac

#DEBHELPER#

exit 0
