# src/cephadm/cephadmlib/templates/agent.service.j2
# generated by cephadm
[Unit]
Description=cephadm agent for cluster {{agent.fsid}}

PartOf=ceph-{{agent.fsid}}.target
Before=ceph-{{agent.fsid}}.target

[Service]
Type=forking
ExecStart=/bin/bash {{agent.daemon_dir}}/unit.run
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=ceph-{{agent.fsid}}.target
