# src/cephadm/cephadmlib/templates/ceph.service.j2
# generated by cephadm
[Unit]
Description=Ceph %i for {{fsid}}

# According to:
#   http://www.freedesktop.org/wiki/Software/systemd/NetworkTarget
# these can be removed once ceph-mon will dynamically change network
# configuration.
After=network-online.target local-fs.target time-sync.target{% if has_docker_engine %} docker.service{% endif %}
Wants=network-online.target local-fs.target time-sync.target
{%- if has_docker_engine %}
Requires=docker.service
{%- endif %}

PartOf=ceph-{{fsid}}.target
Before=ceph-{{fsid}}.target

[Service]
LimitNOFILE=1048576
LimitNPROC=1048576
EnvironmentFile=-/etc/environment
ExecStart=/bin/bash {{ctx.data_dir}}/{{fsid}}/%i/unit.run
ExecStop=-/bin/bash -c 'bash {{ctx.data_dir}}/{{fsid}}/%i/unit.stop'
ExecStopPost=-/bin/bash {{ctx.data_dir}}/{{fsid}}/%i/unit.poststop
KillMode=none
Restart=on-failure
RestartSec=10s
TimeoutStartSec=200
TimeoutStopSec=120
StartLimitInterval=30min
StartLimitBurst=5
{%- if has_podman_engine %}
ExecStartPre=-/bin/rm -f %t/%n-pid %t/%n-cid
ExecStopPost=-/bin/rm -f %t/%n-pid %t/%n-cid
Type=forking
PIDFile=%t/%n-pid
{%- if has_podman_split_version %}
Delegate=yes
{%- endif %}
{%- endif %}

[Install]
WantedBy=ceph-{{fsid}}.target
