# src/cephadm/cephadmlib/templates/sidecar.service.j2
# generated by cephadm
[Unit]
Description=Ceph sidecar %i for {{ sidecar.fsid }}
After=network-online.target local-fs.target time-sync.target
Wants=network-online.target local-fs.target time-sync.target
{%- if has_docker_engine %}
After=docker.service
Requires=docker.service
{%- endif %}
After={{ primary.service_name }}

PartOf=ceph-{{ sidecar.fsid }}.target
Before=ceph-{{ sidecar.fsid }}.target

[Service]
LimitNOFILE=1048576
LimitNPROC=1048576
EnvironmentFile=-/etc/environment
ExecStart=/bin/bash {{ sidecar_script }} start
ExecStop=/bin/bash {{ sidecar_script }} stop
ExecStopPost=-/bin/bash {{ sidecar_script }} poststop
KillMode=none
Restart=on-failure
RestartSec=10s
TimeoutStartSec=200
TimeoutStopSec=120
StartLimitInterval=30min
StartLimitBurst=5
{%- if has_podman_engine %}
ExecStartPre=-/bin/rm -f %t/%n-pid %t/%n-cid
ExecStopPost=-/bin/rm -f %t/%n-pid %t/%n-cid
Type=forking
PIDFile=%t/%n-pid
{%- if has_podman_split_version %}
Delegate=yes
{%- endif %}
{%- endif %}

[Install]
WantedBy=ceph-{{ sidecar.fsid }}.target

